# backend/config/supervisor.medivh.conf
[program:medivh]
command=/usr/bin/example --loglevel=%(ENV_LOGLEVEL)s

